# callastream-export/systemmd/callastream-poller.service
[Unit]
Description=CallaStream Poller (updates state.json from server)
After=network-online.target
Wants=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/callastream
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -u /opt/callastream/player.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
