# callastream-export/systemmd/callastream-kiosk.service
[Unit]
Description=CallaStream Chromium Kiosk
After=systemd-user-sessions.service network-online.target callastream-kiosk-web.service
Wants=network-online.target callastream-kiosk-web.service

[Service]
Type=simple
User=adminjeff
Group=adminjeff
WorkingDirectory=/opt/callastream
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/adminjeff/.Xauthority
Environment=HOME=/home/adminjeff
ExecStart=/usr/bin/xinit /opt/callastream/kiosk-launch.sh -- :0 -nocursor -s 0 -dpms
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
