# callastream-export/systemmd/callastream-kiosk.service
[Unit]
Description=CallaStream WPE WebKit Kiosk (Cog)
After=systemd-user-sessions.service network-online.target callastream-kiosk-web.service
Wants=network-online.target callastream-kiosk-web.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/callastream
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=HOME=/home/pi
Environment=CALLASTREAM_KIOSK_LOG=/var/log/callastream/cog-launch.log
ExecStart=/usr/bin/xinit /opt/callastream/cog-launch.sh -- :0 -nocursor -s 0 -dpms
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
