# callastream-export/systemmd/callastream-player.service.d/logging.conf
[Service]
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal
